# deploy/network-priority.service
[Unit]
Description=Fix WiFi priority over Ethernet
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Удаляем маршрут Ethernet без метрики
ExecStart=/bin/sh -c '/usr/sbin/ip route del default via 192.168.1.1 dev eth0 2>/dev/null || true'
# Добавляем маршрут Ethernet с низким приоритетом (metric 1000)
ExecStart=/usr/sbin/ip route add default via 192.168.1.1 dev eth0 metric 1000
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
